# packages/common/hazelcast-management-center.service
[Unit]
Description=Hazelcast Management Center enables monitoring and management of nodes running Hazelcast.
After=syslog.target network.target

[Service]
User=hazelcast
Group=hazelcast

SuccessExitStatus=143

Type=simple

Environment="JAVA_OPTS=-Dhazelcast.mc.home=/usr/lib/hazelcast-management-center"
WorkingDirectory=/usr/lib/hazelcast-management-center/bin
ExecStart=/usr/lib/hazelcast-management-center/start.sh
ExecStop=/bin/kill -15 $MAINPID

[Install]
WantedBy=multi-user.target
