# packages/deb/hazelcast-management-center/DEBIAN/prerm
#!/bin/bash

if command -v systemctl > /dev/null; then
  systemctl stop hazelcast
  systemctl disable hazelcast
  systemctl daemon-reload
fi

echo "Removing symlinks from /usr/bin"

for FILENAME in /usr/lib/hazelcast-management-center/bin/*mc*; do
  case "${FILENAME}" in
    *bat)
      ;;
    *)
      rm /usr/bin/"$(basename "${FILENAME}")"
      ;;
  esac
done
